# tracker/sunrise-tracker.conf
description "sunrise-tracker"

start on runlevel [345]
stop on runlevel [!345]

respawn

exec /usr/local/bin/sunrise-tracker -log_dir="/var/log"
